# src/fty-shm-cleanup.service.in
[Unit]
Description=fty-shm-cleanup service, one run = one cleanup
After=network.target
# We can have some services using message bus even before EULA
# is accepted; care for their SHM even at this time
PartOf=bios-pre-eula.target

[Service]
# This service runs the cleanup routine and exits
# It is started by timer, so no need to restart automatically
# (in fact systemd throttles that as a loop too quick) nor to
# keep it fake-active as RemainAfterExit=yes.
Type=oneshot
User=bios
ExecStart=@prefix@/bin/fty-shm-cleanup

[Install]
WantedBy=bios-pre-eula.target
