# src/fty-shm-cleanup.timer.in
[Unit]
Description=Timer to regularly trigger the job done by fty-shm-cleanup service
PartOf=bios-pre-eula.target

[Timer]
### Run every night
OnCalendar=*-*-* 04:20:00
# Do not record last-execution times
Persistent=false

[Install]
WantedBy=bios-pre-eula.target
